# systemd/wayland-spice-clipboard.service
[Unit]
Description=Wayland to Spice Clipboard Bridge
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wayland-spice-clipboard
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
